# mailfilter/rootfs/etc/services.d/redis-server/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Redis server...."

exec redis-server /etc/redis.conf