# mailfilter/rootfs/etc/services.d/redis-server/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Mailfilter
# Runs the Redis server
# ==============================================================================

bashio::log.info "Starting Redis server...."

exec redis-server /etc/redis.conf