# tellsticklive/rootfs/etc/services.d/runonce/run
#!/usr/bin/with-contenv bashio

s6-svc -O /var/run/s6/services/runonce

if bashio::config.true "enable_live" && bashio::config.is_empty "live_uuid"; then
    bashio::log.info "UUID has not been entered in config. Generating registration URL"
    bashio::log.info "Follow the link below to register this addon with Telldus Live."
    bashio::log.info "Copy the uuid string in the url and set live_uuid config option."
    exec tellive_core_connector /etc/tellive.conf
fi
