# mail/rootfs/etc/services.d/clamav/run
#!/usr/bin/with-contenv bashio

if bashio::config.true "enable_antivirus"; then
    
    bashio::log.info "Starting Clamd...."
    exec clamd -F
fi
