# mail/rootfs/etc/services.d/clamav/run
#!/usr/bin/with-contenv bashio

if bashio::config.true "enable_antivirus"; then
    
    bashio::log.info "Starting ClamAV antivirus daemon...."
    exec clamd >/dev/null 2>&1
fi
