# tellsticklive/rootfs/etc/services.d/telldusd/run
#!/usr/bin/with-contenv bashio

if bashio::config.true "enable_local"; then
bashio::log.info "Local enabled. Exposing sockets to HomeAssistant..."

socat TCP-LISTEN:50800,reuseaddr,fork UNIX-CONNECT:/tmp/TelldusClient &
socat TCP-LISTEN:50801,reuseaddr,fork UNIX-CONNECT:/tmp/TelldusEvents &
fi

bashio::log.info "Loading Telldus service..."

exec /usr/local/sbin/telldusd --nodaemon
