# mail/rootfs/etc/services.d/clamsmtpd/run
#!/usr/bin/with-contenv bashio
 
bashio::log.info "Starting Clamsmtpd...."
exec clamsmtpd -d 1 -f /etc/clamsmtpd.conf