# mail/rootfs/etc/services.d/dovecot/run
#!/usr/bin/with-contenv bashio
bashio::log.info "Starting Dovecot...."

exec /usr/sbin/dovecot -F -c /etc/dovecot/dovecot.conf
