# mail/rootfs/etc/services.d/postfix/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Postfix...."
exec /usr/libexec/postfix/master -d
