# mailfilter/rootfs/etc/services.d/rspamd/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Rspamd...."

exec rspamd -u rspamd -g rspamd -f