# steve/rootfs/etc/services.d/steve/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: SteVe
# Runs the SteVe OCPP Server
# ==============================================================================

bashio::log.info "Starting SteVe...."
exec java -jar /data/target/steve.jar