# tellsticklive/rootfs/etc/services.d/stdin/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Restarts the STDIN service
# ==============================================================================

bashio::log.info 'STDIN service was stopped, restarting...'
