# mail/rootfs/etc/services.d/syslog/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Syslog...."
exec /sbin/syslogd -n
