# mailfilter/rootfs/etc/services.d/clamav/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Mailfilter
# Runs the Clamav Daemon
# ==============================================================================
   
bashio::log.info "Starting ClamAV antivirus daemon...."
exec clamd >/dev/null 2>&1

